# files/etc/systemd/system/prometheus.service
[Unit]
Description=Prometheus monitoring and alerting 
After=network.target

[Service]
Type=simple
User=prometheus
ExecStart=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --storage.tsdb.retention.time=30d
Restart=on-failure

[Install]
WantedBy=multi-user.target
